# my.service
[Unit]
Wants=my.timer

[Service]
ExecStart=/usr/bin/python3 /root/on_server.py